# l.conf
[program:l]
command=python3 /home/curator2/l/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/l/test.err.log
stdout_logfile=/home/curator2/l/test.out.log